# deploy/daxiao.service
[Unit]
Description=李大霄视频追踪系统
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/daxiao
ExecStart=/opt/daxiao/venv/bin/python -m backend.main serve
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
